# bin/suska_backup
#! /bin/bash
### BEGIN INIT INFO
# Provides:          suska_backup
# Required-Start:    $syslog 
# Required-Stop:     $syslog 
# Should-Start:      
# Should-Stop:
# X-Start-Before:    
# X-Stop-After:      
# Default-Start:     2 3 4 5 6
# Default-Stop:      0
# Short-Description: Einfaches Backup durch Spiegel den User-DAteien
### END INIT INFO

# Install:
# sudo cp -a suska_backup /etc/init.d/
# sudo update-rc.d suska_backup defaults

# Ausführen:
# sudo service suska_backup stop


PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin
#DAEMON=/usr/sbin/cupsd
NAME=suska_backup
#PIDFILE=/var/run/cups/$NAME.pid
DESC="SuskaBackup"

unset TMPDIR

#test -x $DAEMON || exit 0

#mkdir -p /var/run/cups/certs

. /lib/lsb/init-functions

LOGGER="logger -t$NAME "
BACKUP_MOUNTPOINT=/spiegel
ZIEL_PREFIX="$BACKUP_MOUNTPOINT"/

sync_dir() {
    #if [ -n "$success" ] && [ -x /etc/init.d/xprint ]; then
    quelle_eingabe=$1
    quelle=/$quelle_eingabe 
    echo $quelle
    ziel=$ZIEL_PREFIX${quelle_eingabe//\//-}
    $LOGGER "backup sync of $quelle to $ziel in progress ..."
    echo
    echo "######## STARTING BACKUP OF $quelle #########"
    echo
    rsync -av --del --exclude="mtab" --exclude="*/.gvfs"  --exclude="*cache/*" $quelle/ $ziel/ 3>&1 1>&2 2>&3 | tee | $LOGGER
    echo
    echo "######## BACKUP $quelle FINISHED. #########"
    echo
    $LOGGER "backup sync of $quelle done."
    #fi
}



case "$1" in
  start)
	log_begin_msg "Pseudo-Starting $DESC: $NAME"
	success=1
	log_end_msg $?
	;;
  stop)
	log_begin_msg "Stopping for $DESC: $NAME -- now starting real work"
	mount $BACKUP_MOUNTPOINT || $LOGGER "$BACKUP_MOUNTPOINT already mounted??"
	if mountpoint $BACKUP_MOUNTPOINT ; then
	  success=0
	  sync_dir "etc"   || success=4
	  sync_dir "home"  || success=8
	else
	  success=1
	fi
	#log_end_msg $?
	log_end_msg $success
	;;
  reload|force-reload)
       log_begin_msg "Dummy-Reloading $DESC: $NAME"
       log_end_msg $?
       ;;
  restart)
	log_begin_msg "Dummy-Restarting $DESC: $NAME"
	log_end_msg $?
	;;
  status)
	echo -n "Status of $DESC: "
	N=/etc/init.d/${0##*/}
	echo "Usage: $N {start|stop|restart|force-reload|status}" >&2
	exit 1
	;;
esac

exit 0
